# legitdns.timer
[Unit]
Description=legitDNS updater timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
