# legitdns.service.d/example.conf
[Service]
Environment=DDNS_DOMAIN=insert_your_domain_here
Environment=DDNS_TOKEN=insert_your_token_here
